# location-reporter.service
[Unit]
Description=SatPi Location Reporter Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /home/pi/satpi/location-reporter.py daemon
Restart=always
RestartSec=300
User=pi
Group=pi
StandardOutput=journal
StandardError=journal

# Environment variables for email
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=HOME=/home/pi
Environment=PYTHONPATH=/usr/local/lib/python3.9/site-packages
Environment=SATPI_EMAIL_USER=satpi.notifications@outlook.com
EnvironmentFile=-/home/pi/satpi/email-credentials

[Install]
WantedBy=multi-user.target