# satdump-capture.service
[Unit]
Description=SatDump Satellite Capture Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/home/johnc/satpi/satdump-capture.sh auto
Restart=always
RestartSec=30
User=johnc
Group=johnc
StandardOutput=journal
StandardError=journal
WorkingDirectory=/home/johnc/satpi

# Environment
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=HOME=/home/johnc

# Ensure RTL-SDR access
SupplementaryGroups=dialout plugdev

[Install]
WantedBy=multi-user.target