# ssh-upload-daemon.service
[Unit]
Description=SatPi SSH Upload Daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/home/johnc/satpi/ssh-uploader.py daemon
Restart=always
RestartSec=60
User=johnc
Group=johnc
StandardOutput=journal
StandardError=journal
WorkingDirectory=/home/johnc/satpi

# Environment
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=HOME=/home/johnc

[Install]
WantedBy=multi-user.target