# wifi-hunter.service
[Unit]
Description=SatPi WiFi Hunter Service
After=multi-user.target
Wants=network.target

[Service]
Type=simple
ExecStart=/home/pi/satpi/wifi-hunter.sh
Restart=always
RestartSec=10
User=root
StandardOutput=journal
StandardError=journal

# Environment
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

[Install]
WantedBy=multi-user.target